# packages/burnlight/burnlight-0.2.0.tar.gz/burnlight-0.2.0/burnlight/burnlightd.service
[Unit]
Description=burnlightd - A Raspberry Pi GPIO scheduler

[Service]
ExecStart=/usr/bin/python burnlightd
Type=forking

[Install]
WantedBy=multi-user.target