# packages/dynahost/dynahost-1.1.0.tar.gz/dynahost-1.1.0/dynahost.service
# Dynamic host service

[Unit]
Description=Dynamic host service
After=syslog.target

[Service]
Type=simple
User=pi
ExecStart=/usr/bin/dynahost
StandardOutput=syslog
StandardError=syslog

[Install]
WantedBy=multi-user.target
