# packages/cronbackup/cronbackup-1.0.4.tar.gz/cronbackup-1.0.4/cronbackup.service
[Unit]
Description=file backup and synchronize application
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/cronbackup
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
Restart=on-failure
RestartSec=42s

[Install]
WantedBy=multi-user.target
